# backend/stratify.service
# Systemd service file for Stratify AI
# Place this in /etc/systemd/system/stratify.service
# 
# Commands:
#   sudo systemctl start stratify    - Start service
#   sudo systemctl stop stratify     - Stop service
#   sudo systemctl restart stratify  - Restart service
#   sudo systemctl status stratify   - Check status
#   sudo systemctl enable stratify   - Enable on boot

[Unit]
Description=Stratify AI Backend Service
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=notify
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/stratify-ai/backend/app
Environment="PATH=/home/ubuntu/stratify-ai/backend/venv/bin"
Environment="PYTHONPATH=/home/ubuntu/stratify-ai/backend"

# Load environment variables from .env file
EnvironmentFile=/home/ubuntu/stratify-ai/backend/.env

# Start command
ExecStart=/home/ubuntu/stratify-ai/backend/venv/bin/gunicorn \
    -c /home/ubuntu/stratify-ai/backend/gunicorn_conf.py \
    main:app

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stratify-ai

# Security
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
